# skel/.config/autostart/localsend_app.desktop
[Desktop Entry]
Type=Application
Name=localsend_app
Comment=localsend_app startup script
Exec=/usr/share/localsend_app/localsend_app --hidden
StartupNotify=false
Terminal=false
